# packaging/remux-cert-renew.timer
[Unit]
Description=Weekly Tailscale certificate renewal for remux

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
